# conf/collect-logs_nginx.service
[Unit]
Description=Starts Logstash log collecting

[Service]
ExecStart=/usr/share/logstash/bin/logstash -f /etc/logstash/conf.d/logstash_nginx.conf

[Install]
WantedBy=multi-user.target
